# docker/vsftpd/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:vsftpd]
command=vsftpd
