# docker/main/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command = /usr/sbin/php5-fpm --nodaemonize
autostart = true
autorestart = true
