# docker/code/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:init]
command = bash /opt/docker/init.sh
user = root
autostart = false

[program:php-fpm]
command = bash /opt/docker/php-fpm.sh
user = root
autostart = true